# supervisor-himart.conf
# Supervisor 配置文件
# 保存为: /etc/supervisor/conf.d/himart.conf

[program:himart-backend]
command=/home/ubuntu/himart/backend/venv/bin/python manage.py runserver 0.0.0.0:8000
directory=/home/ubuntu/himart/backend
user=ubuntu
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/himart-backend.log
environment=PYTHONUNBUFFERED=1

[program:himart-frontend]
command=/usr/bin/npx serve -s dist -l 3000
directory=/home/ubuntu/himart/frontend
user=ubuntu
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/himart-frontend.log
environment=NODE_ENV="production"

[group:himart]
programs=himart-backend,himart-frontend
priority=999
